# Deployment/local/mysql/crontab/cronfile
0 * * * * /GoSunBigDataDeploy/mysql/sql/sh 24hour_count_table.sh >> /opt/GoSunBigDataDeploy/24hour_count.log 2>&1
0 0 * * * /GoSunBigDataDeploy/mysql/sql/sh device_recogize_table.sh >> /opt/GoSunBigDataDeploy/device_recogize.log 2>&1
0 0 * * * /GoSunBigDataDeploy/mysql/sql/sh imsi_blacklist_table.sh >> /opt/GoSunBigDataDeploy/imsi_blacklist.log 2>&1
0 0 * * * /GoSunBigDataDeploy/mysql/sql/sh fusion_imsi_table.sh >> /opt/GoSunBigDataDeploy/fusion_imsi.log 2>&1
0 0 1 * * /GoSunBigDataDeploy/mysql/sql/sh newpeople_table.sh >> /opt/GoSunBigDataDeploy/newpeople.log 2>&1
0 0 * * * /GoSunBigDataDeploy/mysql/sql/sh outpeople_table.sh >> /opt/GoSunBigDataDeploy/outpeople.log 2>&1